# etc/rc.local
# Randomize system clock to obfuscate boot fingerprints
if [ -x /usr/local/sbin/bootclockrandom ]; then
  /usr/local/sbin/bootclockrandom
fi
